# rootfs/etc/services.d/rundeck/run
#!/usr/bin/with-contenv bash
. /etc/rundeck/profile
exec s6-setuidgid rundeck \
bash -c "exec $rundeckd"
